# docker/cron/crontab
*/1 * * * * /var/www/app/console enhavo:search:index