# docker/cron/crontab
*/1 * * * * /var/www/bin/console enhavo:search:index